# site/init_scripts/runit/ask_app/run
#!/bin/sh

# file with project specific vars
. ./local_vars

exec 2>&1
echo "Start Ask application server..."
exec ${GUNICORN_EXEC} ${GUNICORN_ARGS}

